# etc/crontabs/www-data
# Example cron run every 5 min: */5 * * * * drush cron --root=/var/www/html
*/5 * * * * drush cron --root=/var/www/html
